# files/etc/crontabs/rla
* * * * * /usr/sbin/rla-monitor.sh &
* * * * * /usr/sbin/send-status-monitor.sh &
